# etc/rc.local
#!/bin/sh -e
#
# rc.local - executed at the end of each multiuser runlevel
#
# Make sure that the script will "exit 0" on success or any other
# value on error.

#ln -sf $(find /home/lucas/img/wallpapers/unsplash/ -type f | shuf -n1) /usr/share/backgrounds/default

defrt=`ip route | grep "^default" | head -1`
ip route change $defrt initcwnd 10

cat /etc/hosts.orig > /etc/hosts
echo >> /etc/hosts
echo '# Block List' >> /etc/hosts
cat /etc/hosts.block >> /etc/hosts

exit 0
