# xfce/.config/autostart/keybase_autostart.desktop
[Desktop Entry]
Name=Keybase
Comment=Keybase Filesystem Service and GUI
Type=Application
Exec=env KEYBASE_AUTOSTART=1 run_keybase
Hidden=true

