# etc/systemd/system/lightdm.service.d/override.conf
[Unit]
ConditionPathIsReadWrite=/run/media/lucas/hd_home
