# base/etc/systemd/system/docker.service.d/override.conf
[Unit]
ConditionPathIsReadWrite=/run/media/lucas/hdd

[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H tcp://0.0.0.0:4243 -H unix:///var/run/docker.sock
