# etc/systemd/system/nginx.service.d/restart.conf
[Service]
Restart=always
RemainAfterExit=no
RestartSec=5s

