# gnome/.config/autostart/mega_instances.desktop
[Desktop Entry]
Type=Application
#Exec=/usr/bin/megasync-instances
Name=megasync_instances
Comment=Open all your MEGA instances
