# systemd/etc/systemd/system/media-storage.automount
Description=Automount /media/storage

[Automount]
Where=/media/storage

[Install]
WantedBy=multi-user.target
