# home/modelos/mnt-hdd.automount
[Unit]
Description=HDD Automount

[Automount]
Where=/mnt/hdd

[Install]
WantedBy=multi-user.target
