# etc/systemd/system/mnt-dados.automount
[Unit]
Description=Dados Automount

[Automount]
Where=/mnt/dados

[Install]
WantedBy=multi-user.target
